# deployment/confs/etc/supervisor/conf.d/dashboard.conf
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'

[program:dashboard-server]
directory=/dashboard
command=/run_server.sh
